# files/etc/update-motd.d/99-one-click
#!/bin/sh
#
# Configured as part of the DigitalOcean 1-Click Image build process

myip=$(hostname -I | awk '{print$1}')
cat <<EOF
********************************************************************************

                 **
               **~~**
             **~~~~~~**
           **~~~~~~~~~~**
         **~~~~~~~~~~~~~~**
       **~~~~~~~~~~~~~~~~~~**
     **~~~~~~~~~~~~~~~~~~~~~~**
    **~~~~~~~~~~~~~~~~~~~~~~~~**
   **~~~~~~~~~~~~~~~~~~~~~~~~~~**
  **~~~~~~~~~~~~~~~~~~~~~~~~~~~~**
  **~~~~~~~~~~~~~~~~~~~~~~~~~~~~**
  **~~~~~~~~~~~~~~~~~~~~~++++~~~**
   **~~~~~~~~~~~~~~~~~~~++++~~~**
    ***~~~~~~~~~~~~~~~++++~~~***
      ****~~~~~~~~~~++++~~****
         *****~~~~~~~~~*****
            *************
   
   _       __    ___   ___   ___
  \ \  /  / /\  | |_) / / \ | |_)
   \_\/  /_/--\ |_|   \_\_/ |_| \\
     a web framework for Swift

Welcome to your DigitalOcean Vapor image.

* https://docs.vapor.codes/
* https://swift.org

********************************************************************************
To delete this message of the day: rm -rf $(readlink -f ${0})
EOF
